# apps/voice-gateway-oww/voice-gateway-oww.service
[Unit]
Description=Voice Gateway OWW - Wake Word Detection and Voice Commands
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/code/mqtt-ollama-presentation/apps/voice-gateway-oww
Environment="NODE_ENV=production"
Environment="LOG_LEVEL=info"
Environment="VIRTUAL_ENV=/home/pi/code/mqtt-ollama-presentation/apps/voice-gateway-oww/.venv"
Environment="PATH=/home/pi/code/whisper.cpp/build/bin:/home/pi/.nvm/versions/node/current/bin:/home/pi/code/mqtt-ollama-presentation/apps/voice-gateway-oww/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# OpenWakeWord Configuration
Environment="OWW_MODEL_PATH=models/hey_jarvis_v0.1.onnx"
Environment="OWW_THRESHOLD=0.25"
Environment="OWW_INFERENCE_FRAMEWORK=onnx"

# Audio Configuration
Environment="AUDIO_MIC_DEVICE=plughw:3,0"
Environment="AUDIO_SPEAKER_DEVICE=plughw:2,0"
Environment="AUDIO_SAMPLE_RATE=16000"
Environment="AUDIO_CHANNELS=1"

# Voice Activity Detection (VAD)
Environment="VAD_TRAILING_SILENCE_MS=1500"
Environment="VAD_MAX_UTTERANCE_MS=10000"

# Whisper Speech-to-Text
Environment="WHISPER_MODEL=tiny"
Environment="WHISPER_MODEL_PATH=models/ggml-tiny.bin"

# MQTT Broker
Environment="MQTT_BROKER_URL=mqtt://localhost:1883"
Environment="MQTT_CLIENT_ID=voice-gateway-oww"
Environment="MQTT_USERNAME="
Environment="MQTT_PASSWORD="

# ZWave-JS-UI (for MCP server)
Environment="ZWAVE_UI_URL=http://localhost:8091"
Environment="ZWAVE_UI_AUTH_ENABLED=false"

# Health Check
Environment="HEALTHCHECK_PORT=3002"

# Ollama AI Configuration
Environment="OLLAMA_BASE_URL=http://localhost:11434"
Environment="OLLAMA_MODEL=qwen2.5:0.5b"

# Text-to-Speech (Piper TTS)
Environment="TTS_ENABLED=true"
Environment="TTS_MODEL_PATH=models/piper/en_US-amy-medium.onnx"
Environment="TTS_VOLUME=1.0"
Environment="TTS_SPEED=3.0"

ExecStart=/home/pi/.nvm/versions/node/current/bin/node /home/pi/code/mqtt-ollama-presentation/apps/voice-gateway-oww/src/main.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
